# deploy/systemd/realguard-backup.service
[Unit]
Description=RealGuard application data backup
After=mysql.service realguard-backend.service jianzhen-v2-backend.service
Requires=mysql.service

[Service]
Type=oneshot
User=root
Group=root
UMask=0077
EnvironmentFile=/etc/realguard/realguard-backend.env
EnvironmentFile=/etc/realguard/detector-db.env
EnvironmentFile=/etc/realguard/backup.env
ExecStart=/usr/local/sbin/realguard-backup
TimeoutStartSec=3h
TimeoutStopSec=5m
RuntimeMaxSec=3h
CPUQuota=100%
MemoryHigh=384M
MemoryMax=512M
MemorySwapMax=128M
TasksMax=128
OOMPolicy=stop
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/backups/realguard /opt/realguard-data

[Install]
WantedBy=multi-user.target
