# deploy/systemd/realguard-backend.service
[Unit]
Description=RealGuard Web Backend
After=network.target mysql.service
Requires=mysql.service
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/realguard-server/RealGuard
Environment=REALGUARD_HOST=127.0.0.1
Environment=REALGUARD_PORT=5000
Environment=REALGUARD_SESSION_COOKIE_SECURE=1
Environment=REALGUARD_ENV=production
Environment=REALGUARD_ALERT_WORKER_ENABLED=0
EnvironmentFile=/etc/realguard/session.env
EnvironmentFile=/etc/realguard/realguard-backend.env
EnvironmentFile=/etc/realguard/detector-db.env
EnvironmentFile=-/etc/realguard/sms.env
EnvironmentFile=/etc/realguard/agent.env
EnvironmentFile=/etc/realguard/watermark-precheck.env
ExecStart=/opt/realguard-server/.venv/bin/gunicorn --worker-class gthread --workers 1 --threads 24 --timeout 330 --graceful-timeout 30 --keep-alive 5 --bind 127.0.0.1:5000 --access-logfile - --error-logfile - run:app
Restart=always
RestartSec=3
TimeoutStopSec=45
CPUQuota=150%
MemoryHigh=512M
MemoryMax=640M
MemorySwapMax=256M
TasksMax=256
OOMPolicy=stop
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/opt/realguard-data/uploads /opt/realguard-data /tmp
UMask=0077

[Install]
WantedBy=multi-user.target
