# services/realguard-detection/realguard-detection.service
[Unit]
Description=RealGuard GPU detection service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ymk
Group=ymk
WorkingDirectory=/home/ymk/RealGuard/AIGC_image_detection_system
EnvironmentFile=/etc/realguard/detection-db.env
ExecStart=/home/ymk/miniconda3/envs/realguard/bin/gunicorn --workers 1 --threads 4 --worker-class gthread --bind 0.0.0.0:5000 --timeout 330 --graceful-timeout 45 --keep-alive 5 --access-logfile - --error-logfile - run:app
Restart=always
RestartSec=5
TimeoutStartSec=240
TimeoutStopSec=45
KillSignal=SIGTERM
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
LimitNOFILE=65536
UMask=0027

[Install]
WantedBy=multi-user.target
